# scripts/com.life_dashboard.claude_code_sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.life_dashboard.claude_code_sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/Life-Dashboard/scripts/run_claude_code_sync.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/Life-Dashboard</string>

  <key>StartInterval</key>
  <integer>3600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/life_dashboard_claude_code_sync.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/life_dashboard_claude_code_sync.err.log</string>
</dict>
</plist>
